# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgres-wait]
command=python /app/docker/wait_for_db.py
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api-gateway]
command=bash -c "cd /app/api-gateway/api_gateway && python manage.py migrate --noinput && gunicorn api_gateway.wsgi:application --bind 0.0.0.0:8888 --workers 2"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_DATABASE_URL)s"

[program:customer-service]
command=bash -c "cd /app/customer-service/customer_service && python manage.py migrate --noinput && gunicorn customer_service.wsgi:application --bind 0.0.0.0:8881 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_CUSTOMER_DB_URL)s"

[program:product-service]
command=bash -c "cd /app/product-service/product_service && python manage.py migrate --noinput && gunicorn product_service.wsgi:application --bind 0.0.0.0:8882 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_PRODUCT_DB_URL)s"

[program:cart-service]
command=bash -c "cd /app/cart-service/cart_service && python manage.py migrate --noinput && gunicorn cart_service.wsgi:application --bind 0.0.0.0:8883 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_CART_DB_URL)s",PRODUCT_SERVICE_URL="http://localhost:8882"

[program:catalog-service]
command=bash -c "cd /app/catalog-service/catalog_service && python manage.py migrate --noinput && gunicorn catalog_service.wsgi:application --bind 0.0.0.0:8884 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_CATALOG_DB_URL)s"

[program:order-service]
command=bash -c "cd /app/order-service/order_service && python manage.py migrate --noinput && gunicorn order_service.wsgi:application --bind 0.0.0.0:8885 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_ORDER_DB_URL)s",PAY_SERVICE_URL="http://localhost:8887",SHIP_SERVICE_URL="http://localhost:8886",PRODUCT_SERVICE_URL="http://localhost:8882"

[program:ship-service]
command=bash -c "cd /app/ship-service/ship_service && python manage.py migrate --noinput && gunicorn ship_service.wsgi:application --bind 0.0.0.0:8886 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_SHIP_DB_URL)s"

[program:pay-service]
command=bash -c "cd /app/pay-service/pay_service && python manage.py migrate --noinput && gunicorn pay_service.wsgi:application --bind 0.0.0.0:8887 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_PAY_DB_URL)s"

[program:comment-rate-service]
command=bash -c "cd /app/comment-rate-service/comment_rate_service && python manage.py migrate --noinput && gunicorn comment_rate_service.wsgi:application --bind 0.0.0.0:8888 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_COMMENT_RATE_DB_URL)s"

[program:staff-service]
command=bash -c "cd /app/staff-service/staff_service && python manage.py migrate --noinput && gunicorn staff_service.wsgi:application --bind 0.0.0.0:8889 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_STAFF_DB_URL)s"

[program:manager-service]
command=bash -c "cd /app/manager-service/manager_service && python manage.py migrate --noinput && gunicorn manager_service.wsgi:application --bind 0.0.0.0:8890 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_MANAGER_DB_URL)s"

[program:recommender-ai-service]
command=bash -c "cd /app/recommender-ai-service/recommender_ai_service && python manage.py migrate --noinput && gunicorn recommender_ai_service.wsgi:application --bind 0.0.0.0:8891 --workers 1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="%(ENV_RECOMMENDER_DB_URL)s",PRODUCT_SERVICE_URL="http://localhost:8882",COMMENT_RATE_SERVICE_URL="http://localhost:8888",ORDER_SERVICE_URL="http://localhost:8885",CART_SERVICE_URL="http://localhost:8883"
